# root/etc/services.d/avdc/finish
#!/usr/bin/execlineb -S0

echo "Hello, I'm done."