# root/etc/services.d/avdc/run
#!/usr/bin/execlineb -P

echo "Hello, I'm fine."